# root/etc/services.d/sshfs/finish
#!/command/with-contenv bash

echo "Unmounting ${SSHFS_MNT}"
exec fusermount3 -uz "${SSHFS_MNT}"
