# root/etc/services.d/sshfs/run
#!/command/with-contenv bash

exec s6-setuidgid abc sshfs -f "${SSHFS_HOST}" "${SSHFS_MNT}" ${SSHFS_OPTS}
